# infrastructure/backup/crontab
# =============================================================================
# MyFinPro Backup Cron Configuration
# =============================================================================
# Install with: crontab infrastructure/backup/crontab
# Or add entries to existing crontab: crontab -e
# =============================================================================

# Hourly backup (every hour at minute 0)
0 * * * * /opt/myfinpro/scripts/backup.sh --docker >> /var/log/myfinpro/backup.log 2>&1

# Backup age check every 6 hours (alert if > 2 hours old for hourly frequency)
0 */6 * * * /opt/myfinpro/scripts/check-backup-age.sh --max-age 2 >> /var/log/myfinpro/backup-check.log 2>&1
